# staging_dir/target-aarch64_armv8-a_glibc-2.24/root-layerscape/etc/uci-defaults/luci-i18n-firewall-ko
uci set luci.languages.ko='한국어 (Korean)'; uci commit luci
