# staging_dir/target-aarch64_armv8-a_glibc-2.24/root-layerscape/etc/uci-defaults/luci-i18n-base-sk
uci set luci.languages.sk='SlovenÄina (Slovak)'; uci commit luci
